# etc/network/if-pre-up.d/iptablesload
#!/bin/sh
iptables-restore < /etc/WiTi/iptables.rules
exit 0